# tailscale-autostart
#!/bin/sh /etc/rc.common
# Script tự khởi động tailscaled và tự 'tailscale up' nếu chưa kết nối

START=99
STOP=15

USE_PROCD=1
PROG=/usr/sbin/tailscaled

start_service() {
    echo "[+] Khởi động tailscaled"
    procd_open_instance
    procd_set_param command $PROG
    procd_set_param respawn
    procd_close_instance

    # Đợi daemon lên trước khi kết nối
    sleep 2

    echo "[+] Kiểm tra trạng thái Tailscale..."
    tailscale status >/dev/null 2>&1
    if [ $? -ne 0 ]; then
        echo "[+] Tailscale chưa kết nối, thử khởi tạo lại..."
        tailscale up --ssh
    else
        echo "[✓] Tailscale đã kết nối"
    fi
}

stop_service() {
    echo "[-] Dừng tailscaled"
    killall tailscaled 2>/dev/null
}
